# meta-hpsc-bsp/recipes-core/systemd/files/systemd/system/systemd-hwdb-update.service.d/timeout.conf
[Service]
# Running in Qemu (esp. from NAND) is too slow for default timeout, so increase
TimeoutSec=600s
